# script/crontab
0 * * */5 * sudo /usr/local/bin/ntpclient.hwclock
